# Linux-0.98/Yggdrasil-0.98.3/usr/src/sbin/poeigl/inittab
# inittab for linux, poeigl 1.6
# Format:
# ttyline:termcap-entry:getty-command
tty1:con80x60:/etc/getty 9600 tty1
tty2:con80x60:/etc/getty 9600 tty2
tty3:con80x60:/etc/getty 9600 tty3
tty4:con80x60:/etc/getty 9600 tty4
# tty5:con80x60:/etc/getty 9600 tty5
ttys1:dumb:/etc/getty 9600 tty64
# tty65:dumb:/etc/getty -m -t60 2400 tty65
