# Linux-0.98/Yggdrasil-0.98.3/usr/src/sbin/mount/etc/rc.local
# /etc/rc.local
#
# This file is executed after /etc/rc and consist of local customizations.

# christen ourselves
hostname shiva

# enable screen blanking
setterm -blank 5

# quiet floppy chatter
setfdprm -n /dev/fd0
setfdprm -n /dev/fd1

